# deploy/systemd/cannbench-serve.service
[Unit]
Description=CANNBench service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/cannbench
Environment=PYTHONPATH=/opt/cannbench/src
# To enable GPU JSON upload, replace ExecStart with the commented variant below.
ExecStart=/usr/bin/python3 -m cannbench serve --frontend-dir web/dist --published-dir published --host 127.0.0.1 --port 8000
# ExecStart=/usr/bin/python3 -m cannbench serve --frontend-dir web/dist --published-dir published --host 127.0.0.1 --port 8000 --enable-gpu-upload
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
