# deploy/systemd/cannbench-serve.service
[Unit]
Description=CANNBench service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/cannbench/cannbench-release
Environment=PYTHONPATH=/opt/cannbench/cannbench-release/src
# To enable GPU JSON upload, replace ExecStart with the commented variant below.
ExecStart=/usr/bin/python3 -m cannbench serve --frontend-dir web/dist --published-dir published --host 0.0.0.0 --port 80
# ExecStart=/usr/bin/python3 -m cannbench serve --frontend-dir web/dist --published-dir published --host 0.0.0.0 --port 80 --enable-gpu-upload
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
